# .config/systemd/user/clipmenud.service.d/override.conf
[Service]
Environment="CM_IGNORE_WINDOW=KeePass"